# profile/airootfs/etc/skel/.config/autostart/plasma-firstlogin.desktop
[Desktop Entry]
Type=Application
Name=Plasma First Login Setup
Exec=/usr/local/bin/plasma-firstlogin.sh